# fleet-manager.service
[Unit]
Description=Fleet Control Hub Manager
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/fleet-control-hub-rtl
ExecStart=/usr/bin/node api-server.js
Environment=NODE_ENV=production
Environment=PORT=8080
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=fleet-manager

[Install]
WantedBy=multi-user.target